# config/crontab
* * * * * echo "Hello from cron" >> /var/log/cron.log 2>&1
